# deploy/systemd/jaq.service
[Unit]
Description=JAQ JSON Loader Server (Port 3000)
After=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/jaq/json-loader-server
ExecStart=/home/ubuntu/jaq/target/release/json-loader-server
Restart=always
RestartSec=3
Environment=RUST_LOG=info

[Install]
WantedBy=multi-user.target
